# deployment/etc/cron.d/ckan-harvest
*/3 * * * * root /usr/lib/ckan/bin/python /usr/lib/ckan/bin/paster --plugin=ckanext-harvest harvester run -c /etc/ckan/production.ini > /dev/null
