# deployment/etc/cron.d/ckan-harvest
*/3 * * * * root ckan --plugin=ckanext-harvest harvester run -c /etc/ckan/production.ini > /dev/null
