# deployment/etc/cron.d/ckan-pycsw
# refresh pycsw server keywords weekly
0 0 * * 0 /usr/lib/ckan/bin/paster --plugin=ckanext-spatial ckan-pycsw set_keywords -p /etc/ckan/pycsw.cfg
# sync metadata records daily
0 * * * * /usr/lib/ckan/bin/paster --plugin=ckanext-spatial ckan-pycsw load -p /etc/ckan/pycsw.cfg